# systemd/oroma-orchestrator.service.d/10-restart-policy.conf
# ==============================================================================
# Datei:     /etc/systemd/system/oroma-orchestrator.service.d/10-restart-policy.conf
# Projekt:   ORÓMA v3.8 – Orchestrator (MASTER)
# Stand:     2025-12-27
# Autor:     ORÓMA · KI-JWG-X1
# ==============================================================================
#
# Zweck
# -----
#   Restart-Policy für den Master-Orchestrator.
#
#   Ziel:
#     - Bei Absturz automatisch neu starten.
#     - Maximal 3 Neustarts innerhalb von 5 Minuten.
#     - Danach bleibt der Dienst im "failed"-State, bis manuell eingegriffen wird.
#
# Anpassung
# ---------
#   Diese Werte sind bewusst in einem Drop-In, damit du nur eine Datei ändern musst.
#
#   Nach Änderung:
#     sudo systemctl daemon-reload
#     sudo systemctl restart oroma-orchestrator.service
# ==============================================================================

[Unit]
# 5 Minuten Zeitfenster
StartLimitIntervalSec=300
# max. 3 Neustarts in diesem Zeitfenster
StartLimitBurst=3

[Service]
# Wartezeit bis zum Restart (zusätzlich zur internen Orchestrator-Logik)
RestartSec=3
