# systemd/oroma-orchestrator.service.d/30-asr-whispercpp.conf
[Service]
Environment=OROMA_ASR_BACKEND=whispercpp
Environment=OROMA_WHISPER_ENABLE=false
Environment=OROMA_AUDIO_INPUT_INDEX=0
Environment=OROMA_AUDIO_GAIN=1.0
