# systemd/oroma-orchestrator.service.d/40-ptz-attention.conf
# =============================================================================
# Pfad:      /opt/ai/oroma/systemd/oroma-orchestrator.service.d/40-ptz-attention.conf
# Projekt:   ORÓMA (Offline-Realtime-Organic-Memory-AI)
# Datei:     systemd Drop-In – PTZ Attention / PTZ Motorik Umschaltung
# Version:   v3.8.0+ptz-motor-worker-v1.4d1-live-defaults
# Stand:     2026-05-16
# Autor:     ORÓMA · KI-JWG-X1
# Lizenz:    MIT
# =============================================================================
#
# ZWECK
# ─────
# Diese systemd-Drop-In-Datei steuert die PTZ-Integration im ORÓMA-
# Orchestrator. Bis einschliesslich der alten PTZ-Attention-Phase wurde PTZ
# zyklisch durch den seriellen Orchestrator als One-Shot gestartet:
#
#     python3 -m core.ptz_attention_loop --once
#
# Dieser Pfad ist fuer Diagnose, Scan/Orient/Fixate und langsame Aufmerksamkeit
# geeignet, aber nicht fuer schnelle Motorik bzw. sichtbares PTZ-Following.
# Seit Phase 3a wird die schnelle Motorik durch einen separaten persistenten
# Worker uebernommen:
#
#     tools/ptz_motor_worker.py
#
# WICHTIG
# ───────
# Solange der Motorik-Worker aktiv getestet oder betrieben wird, darf der alte
# Orchestrator-One-Shot nicht parallel laufen. Sonst senden zwei voneinander
# unabhaengige Pfade PTZ-Kommandos an dieselbe Kamera. Das fuehrt zu scheinbar
# zufaelligen Gegenbewegungen, Pendeln und schwer interpretierbaren Logs.
#
# Deshalb ist hier:
#
#     OROMA_ORCH_ENABLE_PTZ_ATTENTION=0
#
# gesetzt. Der alte Code bleibt als Diagnose-/Fallback-Pfad im Projekt erhalten,
# wird aber vom Orchestrator nicht mehr automatisch gestartet.
#
# MOTORIK-WORKER
# ──────────────
# Der neue Motorik-Worker wird in Phase 3a bewusst manuell gestartet:
#
#     sudo -u oroma env \
#       PYTHONPATH=/opt/ai/oroma \
#       OROMA_BASE=/opt/ai/oroma \
#       OROMA_PTZ_DEVICE=/dev/v4l/by-id/usb-EMEET_EMEET_PIXY_A250607001103370-video-index0 \
#       python3 tools/ptz_motor_worker.py --verbose
#
# Live-Kalibrierung 2026-05-10:
# Die unten gepflegten Motorik-Werte sind die bisher beste manuelle Testlinie
# fuer die EMEET PIXY am ORÓMA-Pi: deutlich ruhiger als die ersten 3-Hz-Tests,
# aber weiterhin empfindlich genug fuer Kopf-/Oberkoerperbewegungen. Die Kamera
# zentriert beim Worker-Start auf Pan/Tilt 0/0. Coverage/Arena/Target bleiben
# deaktiviert, damit keine zweite PTZ-Logik parallel Motorbefehle sendet.
#
# Erst nach erfolgreichem Dauerlauf folgt Phase 3b mit eigener systemd-Unit.
#
# =============================================================================

[Service]

# ---------------------------------------------------------------------------
# PTZ Attention One-Shot im Orchestrator deaktivieren
# ---------------------------------------------------------------------------
# 0 = der serielle Orchestrator startet core.ptz_attention_loop --once NICHT.
# 1 = alter One-Shot-Pfad wird wieder durch den Orchestrator ausgefuehrt.
#
# Fuer den neuen PTZ Motor Worker muss dieser Wert auf 0 bleiben.
Environment=OROMA_ORCH_ENABLE_PTZ_ATTENTION=0

# ---------------------------------------------------------------------------
# PTZ Coverage / autonome PTZ-Spieler voruebergehend deaktivieren
# ---------------------------------------------------------------------------
# Diese Jobs bewegen die Kamera autonom fuer Coverage/Policy-Experimente.
# Waehrend der neue PTZ Motor Worker kalibriert wird, muessen sie aus bleiben,
# damit keine zweite Logik die Kamera bewegt. Spaeter koennen sie mit dem neuen
# face-like/upper-body Konzept wieder kontrolliert aktiviert werden.
Environment=OROMA_ORCH_ENABLE_PTZ_COVERAGE_DAILY_RUN=0
Environment=OROMA_ORCH_ENABLE_PTZ_ARENA_DAILY_RUN=0
Environment=OROMA_ORCH_ENABLE_PTZ_TARGET_DAILY_RUN=0

# ---------------------------------------------------------------------------
# Alte Attention-Parameter bleiben als Dokumentation/Fallback erhalten
# ---------------------------------------------------------------------------
# Diese Werte werden vom alten core.ptz_attention_loop genutzt, falls der
# One-Shot-Pfad spaeter zu Diagnosezwecken wieder aktiviert wird.
Environment=OROMA_ORCH_INT_PTZ_ATTENTION=1
Environment=OROMA_ORCH_TIMEOUT_PTZ_ATTENTION=25

Environment=OROMA_PTZ_ATTENTION_AUDIO_GUARD=0
Environment=OROMA_PTZ_ATTENTION_ENSURE_CAM=0

Environment=OROMA_PTZ_ATTENTION_MOTION_LOW=0.008
Environment=OROMA_PTZ_ATTENTION_MOTION_HIGH=0.022
Environment=OROMA_PTZ_THREAT_DEADZONE=0.020
Environment=OROMA_PTZ_THREAT_FORCE_DIST=0.045
Environment=OROMA_PTZ_THREAT_FORCE_ENERGY=0.006

Environment=OROMA_PTZ_ATTENTION_FIXATE_SEC_MIN=2
Environment=OROMA_PTZ_ATTENTION_FIXATE_SEC_MAX=4
Environment=OROMA_PTZ_THREAT_FIXATE_SEC=1
Environment=OROMA_PTZ_ATTENTION_ORIENT_COOLDOWN=4

Environment=OROMA_PTZ_ATTENTION_MOTION_W=96
Environment=OROMA_PTZ_ATTENTION_MOTION_H=54

# ---------------------------------------------------------------------------
# Mapping aus den Live-Tests 2026-05-07
# ---------------------------------------------------------------------------
# Die EMEET PIXY / V4L2-PTZ-Kombination zeigte im Follow-Pfad eine horizontale
# Spiegelung: raw=left musste als mapped=right ausgefuehrt werden.
Environment=OROMA_PTZ_FOLLOW_INVERT_X=1
Environment=OROMA_PTZ_FOLLOW_INVERT_Y=0

# ---------------------------------------------------------------------------
# Startwerte fuer den neuen PTZ Motor Worker
# ---------------------------------------------------------------------------
# Diese ENV-Werte werden nicht vom Orchestrator-One-Shot benoetigt, sind aber
# hier bewusst zentral gepflegt, damit sie bei einem spaeteren systemd-Service
# oder manuellen Start konsistent bleiben.
Environment=OROMA_PTZ_MOTOR_ENABLE=1
Environment=OROMA_PTZ_MOTOR_CENTER_ON_START=1
Environment=OROMA_PTZ_MOTOR_HOME_PAN=0
Environment=OROMA_PTZ_MOTOR_HOME_TILT=0
Environment=OROMA_PTZ_MOTOR_HOME_ZOOM=
Environment=OROMA_PTZ_MOTOR_HOME_SETTLE_SEC=0.70
Environment=OROMA_PTZ_MOTOR_HZ=3
Environment=OROMA_PTZ_MOTOR_MAX_FRAME_AGE=1.50
Environment=OROMA_PTZ_MOTOR_STATE_WRITE_SEC=1.0
Environment=OROMA_PTZ_MOTOR_SUMMARY_SEC=20.0

Environment=OROMA_PTZ_MOTOR_MOTION_W=96
Environment=OROMA_PTZ_MOTOR_MOTION_H=54
Environment=OROMA_PTZ_MOTOR_DEADZONE=0.030
# Live-Tuning 2026-05-10:
# - ENERGY_MIN=0.004 laesst feine Kopfbewegungen durch.
# - STABILITY_MIN=3 beruhigt die Kamera gegen Einzelbild-Zucken.
# - STRONG_SIGNAL_BYPASS=0.400 erlaubt nur wirklich starke Signale ohne Wartezeit.
# - DOWN_CONFIRM_MIN=3 verhindert spontanes Abkippen bei kurzer Bauch-/Armbewegung.
# - UPPER_GAIN/LOWER_DAMPING halten die Aufmerksamkeit bevorzugt im oberen Koerperbereich.
Environment=OROMA_PTZ_MOTOR_ENERGY_MIN=0.004
Environment=OROMA_PTZ_MOTOR_VERTICAL_GAIN=1.35
Environment=OROMA_PTZ_MOTOR_UPPER_BIAS_ENABLE=1
Environment=OROMA_PTZ_MOTOR_UPPER_GAIN=1.35
Environment=OROMA_PTZ_MOTOR_LOWER_DAMPING=0.70
Environment=OROMA_PTZ_MOTOR_HISTORY_N=5
Environment=OROMA_PTZ_MOTOR_DOWN_CONFIRM_MIN=3
Environment=OROMA_PTZ_MOTOR_STABILITY_MIN=3
Environment=OROMA_PTZ_MOTOR_STRONG_SIGNAL_BYPASS=0.400
Environment=OROMA_PTZ_MOTOR_TARGET_ENABLE=1
Environment=OROMA_PTZ_MOTOR_TARGET_DECAY=0.85
Environment=OROMA_PTZ_MOTOR_TARGET_ALPHA=0.45
Environment=OROMA_PTZ_MOTOR_TARGET_CONF_MIN=0.020
Environment=OROMA_PTZ_MOTOR_TARGET_OVERRIDE_RATIO=1.80
Environment=OROMA_PTZ_MOTOR_TARGET_HOLD_TICKS=6
Environment=OROMA_PTZ_MOTOR_TARGET_HOLD_COMMAND=0
Environment=OROMA_PTZ_MOTOR_AXIS_LOCK_ENABLE=1
Environment=OROMA_PTZ_MOTOR_AXIS_LOCK_TICKS=4
Environment=OROMA_PTZ_MOTOR_AXIS_LOCK_OVERRIDE_RATIO=1.65

Environment=OROMA_PTZ_MOTOR_AMOUNT_MIN=2
Environment=OROMA_PTZ_MOTOR_AMOUNT_MAX=3
Environment=OROMA_PTZ_MOTOR_DIST_SCALE=0.18
Environment=OROMA_PTZ_MOTOR_E_SCALE=0.060
Environment=OROMA_PTZ_MOTOR_REVERSAL_GUARD_SEC=0.80
Environment=OROMA_PTZ_MOTOR_REVERSAL_RATIO=1.55