# systemd/99-oroma-usb-noautosuspend.rules
# =============================================================================
# Pfad:    /etc/udev/rules.d/99-oroma-usb-noautosuspend.rules
# Quelle:  /opt/ai/oroma/systemd/99-oroma-usb-noautosuspend.rules
# Projekt: ORÓMA – Headless Lern-KI (Edge)
# Version: v3.7.3
# Stand:   2026-01-30
# Autor:   Jörg + GPT-5.2 Thinking
# =============================================================================
#
# Zweck
# -----
# Deaktiviert USB-Autosuspend (Runtime Power Management) gezielt für ausgewählte
# ORÓMA-Peripherie, um Instabilitäten durch aggressive Suspend/Resume-Zyklen zu
# reduzieren (z.B. UVC-Video-URB Resubmit -19, sporadische Device-Re-Enumeration).
#
# WICHTIGER HINWEIS
# -----------------
# Diese Regel kann KEINE physikalischen "over-current" Ereignisse beheben.
# Over-current Meldungen im Kernel (dmesg) deuten auf Strom-/Hub-/Port-Probleme hin.
# Diese Regel reduziert jedoch eine häufige zusätzliche Ursache für Disconnects:
# Autosuspend/Runtime-PM.
#
# Zielgeräte (IDs)
# ----------------
# 1) EMEET PIXY UVC Kamera
#    idVendor=328f, idProduct=00c0
# 2) Jabra Speak2 55 MS
#    idVendor=0b0e, idProduct=ae6a
#
# Wirkung
# -------
# - setzt power/control auf "on" (keine autosuspend-gesteuerte Runtime-PM)
# - optional: versucht power/autosuspend(-delay) zu neutralisieren, wenn Attribute existieren
#
# Installation (durch User)
# -------------------------
# Diese Datei liegt im ORÓMA-Repo unter /opt/ai/oroma/systemd/.
# Kopiere sie nach /etc/udev/rules.d/ und lade udev neu:
#   sudo udevadm control --reload-rules && sudo udevadm trigger
#
# Diagnose
# --------
# Prüfen:
#   cat /sys/bus/usb/devices/*/idVendor | ...
#   cat /sys/bus/usb/devices/<dev>/power/control
#
# =============================================================================

ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="328f", ATTR{idProduct}=="00c0", TEST=="power/control", ATTR{power/control}="on"
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0b0e", ATTR{idProduct}=="ae6a", TEST=="power/control", ATTR{power/control}="on"

# Best-effort: falls das System power/autosuspend(-delay) anbietet, neutralisieren.
# Einige Kernel/udev Kombinationen akzeptieren negative Werte nicht per ATTR.
# Daher hier nur die sichere Einstellung power/control=on oben – und zusätzlich
# ein optionaler RUN-Hook, der robust per shell schreibt (wenn Attribute existieren).

ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="328f", ATTR{idProduct}=="00c0", RUN+="/bin/sh -c 'test -e /sys$devpath/power/autosuspend && echo -1 > /sys$devpath/power/autosuspend 2>/dev/null || true; test -e /sys$devpath/power/autosuspend_delay_ms && echo -1 > /sys$devpath/power/autosuspend_delay_ms 2>/dev/null || true'"

ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="0b0e", ATTR{idProduct}=="ae6a", RUN+="/bin/sh -c 'test -e /sys$devpath/power/autosuspend && echo -1 > /sys$devpath/power/autosuspend 2>/dev/null || true; test -e /sys$devpath/power/autosuspend_delay_ms && echo -1 > /sys$devpath/power/autosuspend_delay_ms 2>/dev/null || true'"
