# systemd/oroma-exportgate.service
# =============================================================================
# Datei:    /etc/systemd/system/oroma-exportgate.service
# Projekt:  ORÓMA v3.8 – Export-Gate (Delay/Qualität + TAR-Export)
# Stand:    2025-11-30
# Autor:    ORÓMA · KI-JWG-X1
# =============================================================================

[Unit]
Description=ORÓMA v3.8 – Export-Gate (Registry → Bundle-Export)
After=network-online.target oroma.service
Wants=network-online.target
ConditionPathExists=!/opt/ai/oroma/.use_orchestrator

[Service]
Type=oneshot
WorkingDirectory=/opt/ai/oroma
SyslogIdentifier=oroma-exportgate

User=oroma
Group=oroma
UMask=002

EnvironmentFile=-/opt/ai/oroma/.env.systemd
EnvironmentFile=-/opt/ai/oroma/.env

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=UTF-8
Environment=TZ=Europe/Berlin
Environment=PYTHONPATH=/opt/ai/oroma
Environment=OROMA_BASE=/opt/ai/oroma
Environment=PATH=/usr/local/bin:/usr/bin:/bin

ExecStartPre=/usr/bin/env bash -lc 'mkdir -p /opt/ai/oroma/exports /opt/ai/oroma/logs'
ExecStart=/usr/bin/python3 -m core.export_gate

StandardOutput=append:/opt/ai/oroma/logs/exportgate.out.log
StandardError=append:/opt/ai/oroma/logs/exportgate.err.log

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
ProtectClock=yes
RestrictRealtime=yes
ReadWritePaths=/opt/ai/oroma/logs /opt/ai/oroma/exports /opt/ai/oroma/data

[Install]
WantedBy=multi-user.target
