# systemd/oroma-ptz-motor-reward-collector.service
# =============================================================================
# Datei:    /etc/systemd/system/oroma-ptz-motor-reward-collector.service
# Projekt:  ORÓMA (Offline-Realtime-Organic-Memory-AI)
# Modul:    PTZ Motor Reward Collector – persistenter Utility-Slow-Loop
# Version:  v3.8.0+ptz-motor-reward-collector-systemd-v1.0
# Stand:    2026-05-20
# Autor:    ORÓMA · ChatGPT Patch-Gate
# Lizenz:   MIT
# =============================================================================
#
# ZWECK
# ─────
# Diese Unit startet den PTZ Motor Reward Collector als eigenen systemd-Dienst.
# Der Collector ist bewusst vom schnellen PTZ Motor Worker getrennt und liest
# nur den vom Worker geschriebenen State unter data/state/ptz_motor_state.json.
# Daraus erzeugt er langsame Utility-/Reward-Signale für rewards_log und die
# spätere Dream-/Policy-Verdichtung.
#
# WICHTIGE ENTSCHEIDUNG
# ─────────────────────
# Der Collector wird NICHT in den Worker-Hot-Loop eingebaut und NICHT dem
# Orchestrator überlassen. Er bleibt ein separater, dauerhaft laufender
# Slow-Loop und wird systemd-seitig an den Worker gekoppelt.
#
#   - Worker startet -> Collector wird mit gestartet
#   - Worker stoppt/restartet -> Collector stoppt/restartet mit
#   - Keine DB-Schreibzugriffe im Worker-Hot-Path
#   - Keine manuelle Zusatz-Session mehr nötig
#
# LOGS
# ────
#   /opt/ai/oroma/logs/ptz_motor_reward_collector.out.log
#   /opt/ai/oroma/logs/ptz_motor_reward_collector.err.log
#
# =============================================================================

[Unit]
Description=ORÓMA – PTZ Motor Reward Collector
After=oroma.service oroma-ptz-motor-worker.service
Wants=oroma.service
Requires=oroma-ptz-motor-worker.service
PartOf=oroma-ptz-motor-worker.service
ConditionPathExists=/opt/ai/oroma/.use_orchestrator

[Service]
Type=simple
User=oroma
Group=oroma
UMask=002
WorkingDirectory=/opt/ai/oroma
SyslogIdentifier=oroma-ptz-motor-reward-collector

EnvironmentFile=-/opt/ai/oroma/.env.systemd
EnvironmentFile=-/opt/ai/oroma/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONIOENCODING=UTF-8
Environment=TZ=Europe/Berlin
Environment=PYTHONPATH=/opt/ai/oroma
Environment=OROMA_BASE=/opt/ai/oroma
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=OROMA_LOG_DIR=/opt/ai/oroma/logs
Environment=OROMA_DBW_ENABLE=1
Environment=OROMA_PTZ_MOTOR_COLLECTOR_ENABLE=1

ExecStartPre=+/usr/bin/env bash -lc 'mkdir -p /opt/ai/oroma/logs /opt/ai/oroma/data/state && touch /opt/ai/oroma/logs/ptz_motor_reward_collector.out.log /opt/ai/oroma/logs/ptz_motor_reward_collector.err.log && chown oroma:oroma /opt/ai/oroma/logs/ptz_motor_reward_collector.out.log /opt/ai/oroma/logs/ptz_motor_reward_collector.err.log && chown oroma:oroma /opt/ai/oroma/data/state'
ExecStart=/usr/bin/python3 /opt/ai/oroma/tools/ptz_motor_reward_collector.py

Restart=on-failure
RestartSec=2
TimeoutStartSec=30
TimeoutStopSec=10
KillSignal=SIGTERM
KillMode=control-group

StandardOutput=append:/opt/ai/oroma/logs/ptz_motor_reward_collector.out.log
StandardError=append:/opt/ai/oroma/logs/ptz_motor_reward_collector.err.log

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
ReadWritePaths=/opt/ai/oroma/logs /opt/ai/oroma/data /opt/ai/oroma/exports /opt/ai/oroma/uploads /opt/ai/oroma/archives

[Install]
WantedBy=multi-user.target
