# systemd/oroma.service.d/10-env.conf
[Service]
EnvironmentFile=
EnvironmentFile=-/opt/ai/oroma/.env.systemd
EnvironmentFile=-/opt/ai/oroma/.env