# systemd/oroma.service.d/12-groups-audio.conf
[Service]
SupplementaryGroups=video render input gpio i2c audio
