# systemd/oroma.service.d/20-vision.conf
# /etc/systemd/system/oroma.service.d/20-vision.conf
[Service]
Environment=OROMA_CAMERA_BACKEND=opencv
Environment=VISION_BACKEND=opencv
Environment=OROMA_VISION_BACKEND=opencv
Environment=OROMA_VISION_DEVICE=/dev/v4l/by-id/usb-EMEET_EMEET_PIXY_A250607001103370-video-index0
Environment=OROMA_VISION_W=1280
Environment=OROMA_VISION_H=720
Environment=OROMA_VISION_FPS=15
Environment=OROMA_VISION_FOURCC=MJPG
Environment=OROMA_VISION_BUFFERSIZE=2
Environment=OROMA_AV_SNAPS=1
Environment=OROMA_AV_SNAPS_EVERY_TICKS=20
Environment=OROMA_AV_MIN_MOTION=0.0
Environment=OROMA_AV_MIN_QUALITY=0.0
Environment=OROMA_VISION_INFER=1
Environment=OROMA_VISION_INFER_EVERY_TICKS=8
Environment=OROMA_CAM_TRAIN_K=5
Environment=OROMA_CAM_TRAIN_MIN_Q=0.0
Environment=OROMA_CAM_TRAIN_WINDOW_SEC=86400
Environment=OROMA_MODELS_DIR=/opt/ai/oroma/models
# Falls du Luminanz über Kamera nutzt: (zum Testen ggf. 'off')
Environment=OROMA_LIGHT_SOURCE=camera

# Kamera-Ausrichtung:
#   0   = keine Rotation
#   180 = Kopf-ueber / umgedreht (vflip+hflip)
Environment=OROMA_VISION_ROTATE_DEG=0