# systemd/oroma.service.d/25-vision-rotate.conf
# ==============================================================================
# Datei:     systemd/oroma.service.d/25-vision-rotate.conf
# Projekt:   ORÓMA – Vision / Camera Orientation (DeviceHub software rotate)
# Stand:     2025-12-28
# Autor:     ORÓMA · KI-JWG-X1 (Jörg) + GPT-5.2 Thinking
# ==============================================================================
# Zweck
# -----
#   Dreht das Kamerabild *softwareseitig* direkt im DeviceHub, bevor es in die
#   Vision-Pipeline (Tokens, MJPEG, UI) geht.
#
#   Motivation:
#     - Manche PiCam-Mounts liefern ein "auf dem Kopf" stehendes Bild.
#     - 180° ist der typische Fix.
#
# Nutzung
# ------
#   0/90/180/270 (Grad, CCW).
#   180 ist in der Praxis meist "kopfstehend" -> richtig.
#
# Hinweis:
#   - 180° wird als schneller Flip umgesetzt (sehr geringe CPU-Last).
#   - 90/270 nutzt np.rot90 (ebenfalls ok, aber etwas teurer).
# ==============================================================================

[Service]
Environment=OROMA_VISION_ROTATE=180
