# systemd/oroma.service.d/30-ptz.conf
# /etc/systemd/system/oroma.service.d/30-ptz.conf
# =============================================================================
# Pfad:    /etc/systemd/system/oroma.service.d/30-ptz.conf
# Projekt: ORÓMA – Headless Lern-KI (Edge)
# Version: v3.7.3
# Stand:   2026-01-26
# Autor:  GPT-5.2 Thinking
# =============================================================================
#
# Zweck
# -----
# Dieses systemd Drop-In aktiviert die optionale PTZ-Steuerung (Pan/Tilt/Zoom)
# für eine per USB/UVC angebundene Kamera, die V4L2-PTZ-Controls bereitstellt.
# In deinem Setup ist das die EMEET PIXY auf /dev/video8, welche u. a. folgende
# Controls exportiert:
#   - pan_absolute  (z. B. -540000 .. 540000, step=3600)
#   - tilt_absolute (z. B. -324000 .. 324000, step=3600)
#   - zoom_absolute (z. B. 100 .. 150)
#
# Hintergrund (warum nötig)
# -------------------------
# Die PTZ-Implementierung in ORÓMA ist bewusst "fail-closed".
# Das heißt: wenn OROMA_PTZ_DEVICE nicht gesetzt ist (oder v4l2-ctl fehlt),
# wird PTZ als "disabled" gemeldet und in der UI erscheint "Device: n/a".
# Bei der Umstellung des Video-Backends auf OpenCV/V4L2 (PIXY als Quelle) wurde
# die PTZ-ENV nicht mehr gesetzt – daher war PTZ nicht mehr erkennbar.
#
# Wirkung / Priorität
# -------------------
# Drop-In Dateien in /etc/systemd/system/oroma.service.d/ werden vom systemd
# Dienst oroma.service zusätzlich geladen.
# Dieses Drop-In setzt nur PTZ-spezifische Variablen und greift nicht in
# Kamera/Vision-Parameter ein.
#
# Deploy / Aktivierung
# --------------------
# 1) Datei nach Live-System kopieren:
#      /opt/ai/oroma/systemd/oroma.service.d/30-ptz.conf  →  /etc/systemd/system/oroma.service.d/30-ptz.conf
#    (oder direkt in /etc/... ablegen, wenn du systemd-Units so pflegst)
# 2) systemd neu laden und Dienste neu starten:
#      sudo systemctl daemon-reload
#      sudo systemctl restart oroma.service oroma-orchestrator.service
# 3) Wirksamkeit prüfen:
#      systemctl show oroma.service -p Environment | tr ' ' '\n' | egrep 'OROMA_PTZ_'
#
# Wichtige Hinweise (Stabilität)
# ------------------------------
# - PTZ läuft unabhängig vom Video-Capture (OpenCV/picamera2). Es schreibt nur
#   V4L2-Controls auf das angegebene Device.
# - Cooldown verhindert "Move-Spam" und stabilisiert Bild + USB + UI.
# - Wenn /dev/video8 sich ändert (anderer Port), passe OROMA_PTZ_DEVICE an.
#

[Service]
# PTZ-Device der EMEET PIXY (V4L2-Controls sind auf /dev/video8 vorhanden)
Environment=OROMA_PTZ_DEVICE=/dev/v4l/by-id/usb-EMEET_EMEET_PIXY_A250607001103370-video-index0

# Rate-Limit für PTZ-Kommandos (ms). 250ms ist robust und verhindert Spikes.
Environment=OROMA_PTZ_COOLDOWN_MS=250
