# systemd/oroma.service.d/35-audio.conf
[Service]
# --- ORÓMA Audio – Jabra Speak2 55 MS (Capture + Playback) ---
# WICHTIG: systemd 'Environment=' trennt an Spaces, daher bei Gerätenamen IMMER quoten:
#   Environment="OROMA_AUDIO_OUTPUT_NAME=Jabra Speak2 55 MS"
# Sonst landet nur "Jabra" in der Env und das Device-Matching kann fehlschlagen.

Environment=OROMA_AUDIO_ENABLE=true

# Input ist bei dir sicher Index 0 (Jabra)
#Environment=OROMA_AUDIO_INPUT_INDEX=0
#Environment="OROMA_AUDIO_INPUT_NAME=Jabra Speak2 55 MS"

# Playback ebenfalls über Jabra
#Environment="OROMA_AUDIO_OUTPUT_NAME=Jabra Speak2 55 MS"

# ORÓMA intern 16kHz mono (Capture)
#Environment=OROMA_AUDIO_SR=16000
#Environment=OROMA_AUDIO_CH=1

# Gain in DeviceHub ist dB: 0.0 = neutral
Environment=OROMA_AUDIO_GAIN=0.0

# Optional: Mic dauerhaft aktiv halten (Boot)
Environment=OROMA_AUDIO_ALWAYS_ON=1

# Audio-SnapTokens (leichte Features in DB, kein WAV) – analog zu Kamera-AV-SnapTokens
Environment=OROMA_AUDIO_SNAPS=1
Environment=OROMA_AUDIO_SNAPS_EVERY_TICKS=20
Environment=OROMA_AUDIO_SNAPS_WINDOW_SEC=0.50
# 1 = nur wenn VAD anspringt (Sprache/aktiv), 0 = auch Raum-/Stille-Tokens
Environment=OROMA_AUDIO_SNAPS_SPEECH_ONLY=0
Environment=OROMA_AUDIO_SNAPS_MIN_RMS=0.006
