# systemd/oroma.service.d/45-nmr-lite.conf
# =============================================================================
# Pfad:    systemd/oroma.service.d/45-nmr-lite.conf
# Projekt: ORÓMA – NMR-Lite Live Path / Observability Gate
# Stand:   2026-05-25
# Zweck:
#   Aktiviert den leichten NMR-Lite-Livepfad im oroma.service-Prozess.
#
# Hintergrund:
#   core.agent_loop registriert den NMR-Lite-Hook bewusst nur, wenn
#   OROMA_NMR_ENABLE=1 gesetzt ist. Der NMR-Lite-Core selbst kann importierbar
#   und testbar sein, ohne dass der produktive AgentLoop ihn tickt. Dieses
#   Drop-In macht die Aktivierung im echten systemd-Lauf eindeutig sichtbar.
#
# Wirkung:
#   - OROMA_NMR_ENABLE=1: aktiviert NMR-Lite global für den Service-Prozess.
#   - OROMA_NMR_AGENTLOOP=1: registriert/tickt den Hook im AgentLoop.
#   - OROMA_NMR_PERSIST=1: erlaubt gedrosselte Metrik-Persistenz.
#   - OROMA_NMR_PERSIST_WINDOW_SEC=30: verhindert Metrik-Spam.
#   - OROMA_NMR_METRIC_PREFIX=nmr: erzeugt Serien wie nmr:pe/nmr:pe_ema.
#
# Betriebsregel:
#   Der Pfad ist best-effort. Fehler im NMR-Lite-Tick werden im Code gedrosselt
#   geloggt und dürfen den AgentLoop nicht destabilisieren.
# =============================================================================

[Service]
Environment=OROMA_NMR_ENABLE=1
Environment=OROMA_NMR_AGENTLOOP=1
Environment=OROMA_NMR_PERSIST=1
Environment=OROMA_NMR_PERSIST_WINDOW_SEC=30
Environment=OROMA_NMR_METRIC_PREFIX=nmr
