# systemd/oroma.service.d/90-orchestrator-master.conf
# =============================================================================
# Pfad:    systemd/oroma.service.d/90-orchestrator-master.conf
# Projekt: ORÓMA – Orchestrator Master-Bindings (systemd Drop-In)
# Stand:   2025-12-27
# Zweck:
#   Macht 'oroma-orchestrator.service' zum Master:
#     - Stop des Orchestrators -> stoppt oroma.service (und umgekehrt)
#     - Start von oroma.service zieht den Orchestrator mit hoch (sicherer Default)
#
# Hinweis:
#   In deinem Setup existiert die Datei /opt/ai/oroma/.use_orchestrator.
#   Diese Drop-In greift NUR in Orchestrator-Modus.
# =============================================================================

[Unit]
ConditionPathExists=/opt/ai/oroma/.use_orchestrator

# Wenn jemand aus Versehen "systemctl start oroma" macht, soll trotzdem der
# Orchestrator laufen (Master-Workflow).
Wants=oroma-orchestrator.service

# Harte Kopplung: Wenn der Orchestrator endet, endet oroma ebenfalls.
BindsTo=oroma-orchestrator.service
PartOf=oroma-orchestrator.service
